# build/darwin/com.wireguide.helper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>io.github.steiale.wireguide-plus.helper</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Library/PrivilegedHelperTools/io.github.steiale.wireguide-plus.helper</string>
        <string>--helper</string>
        <string>--socket=/var/run/wireguide-plus/wireguide-plus.sock</string>
        <string>--uid=__UID__</string>
        <string>--data-dir=/Library/Application Support/wireguide-plus</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardErrorPath</key>
    <string>/var/log/wireguide-plus-helper.log</string>
    <key>StandardOutPath</key>
    <string>/var/log/wireguide-plus-helper.log</string>
</dict>
</plist>
